# XFCE/files/.config/autostart/onboard.desktop
[Desktop Entry]
Type=Application
Name=Onboard
Exec=onboard
X-GNOME-Autostart-enabled=true