# LXDE/files/.config/autostart/qjoypad.desktop
[Desktop Entry]
Type=Application
Name=QJoyPad
Exec=qjoypad
X-GNOME-Autostart-enabled=true